# templates/crossbar.conf.erb
description "Crossbar.io"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 20 5

setuid <%= @user %>
setgid <%= @user %>

exec /opt/crossbar/bin/crossbar start --cbdir=/home/<%= @user %>/.crossbar --loglevel=<%= @log_level %> --logtofile --logdir /home/<%= @user %>